# deploy/systemd/queuebus-gbis-collector.service
[Unit]
Description=QueueBus GBIS collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=queuebus
WorkingDirectory=/home/queuebus/queue-bus
ExecStart=/usr/bin/npm run gbis:collect
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
